# net-proxy/mtproxy/files/mtproxy-initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

: ${MTPD_USER:=mtproxy}
: ${MTPD_GROUP:=mtproxy}

: ${MTPD_CONFIG:=${RC_PREFIX%/}/etc/mtproxy.conf}
: ${MTPD_PIDFILE:=${RC_PREFIX%/}/run/${SVCNAME}.pid}

: ${MTPD_DIR:=${RC_PREFIX%/}/run/mtproxy}
: ${MTPD_LOG:=${RC_PREFIX%/}/var/log/mtproxy/mtproxy.log}
: ${MTPD_RUN:=${RC_PREFIX%/}/usr/bin/mtproxy.sh}

description="Simple MT-Proto proxy"

command="${MTPD_RUN}"
command_args="run"
command_background="yes"
command_user="${MTPD_USER}:${MTPD_GROUP}"

directory="${MTPD_DIR}"
output_log="${MTPD_LOG}"
pidfile="${MTPD_PIDFILE}"

start_stop_daemon_args="--env CONFIG=${MTPD_CONFIG} ${MTPD_SSD_ARGS}"

depend() {
	after ntp-client
	use net
}

start_pre() {
	checkpath -d -m 0755 -o "${MTPD_USER}:${MTPD_GROUP}" /run/mtproxy
}
